# plane/rootfs/etc/s6-overlay/s6-rc.d/valkey/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Valkey (Redis-compatible)
# ==============================================================================

bashio::log.info "Starting Valkey..."

mkdir -p /data/valkey

exec valkey-server /etc/valkey/valkey.conf
