# plane/rootfs/etc/s6-overlay/s6-rc.d/postgres/finish
#!/command/with-contenv bashio
bashio::log.info "PostgreSQL stopped."
